# salt/appliance/metric/pghero/pghero-query-stats.timer
[Unit]
Description=pghero:capture_query_stats

[Timer]
# every 5 minutes
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
